# ops/deploy/backup.cron.example
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Backup diário do MySQL às 02:30
30 2 * * * cd /home/Mágico && ./ops/backup/mysql-backup.sh >> /var/log/sonhomagico-backup.log 2>&1

# Auditoria operacional semanal aos domingos às 03:30
30 3 * * 0 cd /home/Mágico && ./ops/security/security-audit.sh >> /var/log/sonhomagico-security.log 2>&1
